# glftpd.cron
#!/bin/sh
/home/services/glftpd/bin/reset -r /etc/glftpd.conf
